# fs/etc/supervisor.d/magentoQueueConsumer_inventoryQtyCounter.ini
[program:magentoQueueConsumer_inventoryQtyCounter]
command = /usr/bin/php /var/www/dist/bin/magento queue:consumers:start inventoryQtyCounter
user = root
autostart = false
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
