# fs/etc/supervisor.d/magentoQueueConsumer_sharedCatalogUpdateCategoryPermissions.ini
[program:magentoQueueConsumer_sharedCatalogUpdateCategoryPermissions]
command = /usr/bin/php /var/www/dist/bin/magento queue:consumers:start sharedCatalogUpdateCategoryPermissions
user = root
autostart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
